# initscripts/ovirt-ha-broker.service.in
[Unit]
Description=oVirt Hosted Engine High Availability Communications Broker
# See https://bugzilla.redhat.com/1164226
ConditionFileNotEmpty=/etc/ovirt-hosted-engine/hosted-engine.conf

# Make sure the ovirt-node persistence service restored
# all config files before starting this on Node
After=ovirt-early.service

[Service]
EnvironmentFile=-/etc/sysconfig/ovirt-ha-broker
ExecStart=@ENGINE_HA_BINDIR@/ovirt-ha-broker
Restart=on-failure
StandardOutput=null
User=@VDSM_USER@
Group=@VDSM_GROUP@
RuntimeDirectory=@PACKAGE_NAME@
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
