# initscripts/ovirt-ha-agent.service.in
[Unit]
Description=oVirt Hosted Engine High Availability Monitoring Agent
Wants=ovirt-ha-broker.service
Wants=vdsmd.service
Wants=sanlock.service
After=ovirt-ha-broker.service
After=vdsmd.service
After=sanlock.service

# Make sure the ovirt-node persistence service restored
# all config files before starting this on Node
After=ovirt-early.service

# See https://bugzilla.redhat.com/1164226
ConditionFileNotEmpty=/etc/ovirt-hosted-engine/hosted-engine.conf

[Service]
EnvironmentFile=-/etc/sysconfig/ovirt-ha-agent
ExecStart=@ENGINE_HA_BINDIR@/ovirt-ha-agent
RestartSec=10s
Restart=always
StandardOutput=null
User=@VDSM_USER@
Group=@VDSM_GROUP@

[Install]
WantedBy=multi-user.target
Also=ovirt-ha-broker.service
